# config/systemd/user/good-night.service
[Unit]
Description=Good Night service

[Service]
Type=simple
ExecStart=/usr/bin/notify-send -i gnome-info "ねろ"

[Install]
WantedBy=default.target
